# docker/supervisor.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:etherpad]
directory=/opt/etherpad-lite
command=node node_modules/ep_etherpad-lite/node/server.js
user=root
autostart=true
autorestart=true

